# workspace/config/supervisor/conf.d/user-apps.conf

# workspace/config/supervisor/conf.d/user-apps.conf
[program:user-apps]
command=/workspace/scripts/run-app.sh %(program_name)s
directory=/workspace_data/%(program_name)s
user=www-data
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/%(program_name)s.err.log
stdout_logfile=/var/log/supervisor/%(program_name)s.out.log